# config/wbmachine-process-pending-archives.timer
[Unit]
Description=Process pending archives every minute

[Timer]
OnCalendar=*-*-* *:*:00
